# files/supervisord.ini
[supervisord]
nodaemon=false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log

[inet_http_server]
port = *:9001

[program:slapd]
priority=20
directory=/tmp
command=/sbin/slapd -u ldap -h "ldap:/// " -d 1
user=root
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/slapd.log
stderr_logfile=/var/log/supervisor/slapd.log
